# supervisord.conf
[supervisord]
nodaemon=true

[program:chatscript]
command=/entrypoint-chatscript.sh -D

[program:slackbot]
command=/entrypoint-slackbot.sh -D
